# internal/render/testdata/gl-mt6000/expected/etc/crontabs/root
0 */6 * * * /usr/bin/update-antifilter.sh
0 4 * * 0 /usr/bin/update-geosite.sh
